# infra/vm/quadlet/web.container
# Event Copilot web (Next.js standalone) on node:22-slim. The standalone bundle
# is bind-mounted from /srv/event-copilot/app (rsynced by CI); we run the generated
# server.js (never `require('next')`). Published to host loopback only; the host
# nginx terminates TLS and reverse-proxies :443 → 127.0.0.1:3000.
[Unit]
Description=Event Copilot — web (Next.js standalone)
After=ecnet-network.service chroma.service d1shim.service
Requires=ecnet-network.service
Wants=chroma.service d1shim.service

[Container]
ContainerName=web
Image=docker.io/library/node:22-slim
PublishPort=127.0.0.1:3000:3000
# The packaged standalone root (contains apps/web/server.js, node_modules, packages).
Volume=/srv/event-copilot/app:/app:ro,Z
WorkingDir=/app
Network=ecnet.network
EnvironmentFile=/srv/event-copilot/shared/web.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Exec=node apps/web/server.js

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
