# infra/vm/quadlet/web.container
# Event Copilot web (Next.js standalone) on node:22-slim. The standalone bundle
# is bind-mounted from /srv/event-copilot/app (rsynced by CI); we run the generated
# server.js (never `require('next')`). Published to host loopback only; the host
# nginx terminates TLS and reverse-proxies :443 → 127.0.0.1:3000.
[Unit]
Description=Event Copilot — web (Next.js standalone)
After=ecnet-network.service chroma.service d1shim.service snapshot-data-volume.service
# Requires only the network (hard dep); chroma/d1shim are Wants (soft) so a
# transient data-container outage doesn't take the whole web app down — the
# routes that don't need them keep serving, and RAG/reminders recover when they
# come back.
Requires=ecnet-network.service
Wants=chroma.service d1shim.service snapshot-data-volume.service

[Container]
ContainerName=web
# Pinned to digest (was node:22-slim) so a surprise upstream re-pull can't shift
# runtime behavior during the demo window.
Image=docker.io/library/node:22-slim@sha256:d9f850096136edbc402debdd8729579a288aac64574ada0ff4db26b6ae58b0b2
PublishPort=127.0.0.1:3000:3000
# The packaged standalone root (contains apps/web/server.js, node_modules, packages).
Volume=/srv/event-copilot/app:/app:ro,Z
# Runtime snapshot volume (written by the crawl→re-ingest loop). data.ts hot-reloads
# this file (mtime-gated) so the Schedule/now/next UI updates without a redeploy;
# falls back to the bundled snapshot if the file is missing/invalid.
Volume=snapshot-data.volume:/snap:ro
WorkingDir=/app
Network=ecnet.network
EnvironmentFile=/srv/event-copilot/shared/web.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=SNAPSHOT_PATH=/snap/snapshot.json
Exec=node apps/web/server.js

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
